# fntx-cli/ibkr/alm_daily_import.timer
[Unit]
Description=Run ALM Daily Import at 6:30 AM EST
Requires=alm_daily_import.service

[Timer]
# Run daily at 6:30 AM EST (11:30 UTC)
OnCalendar=*-*-* 11:30:00
Persistent=true

# Randomize by up to 5 minutes to avoid API rate limits
RandomizedDelaySec=300

[Install]
WantedBy=timers.target