# backend/services/alm-automation.service
[Unit]
Description=ALM Automation Service
After=network.target postgresql.service

[Service]
Type=oneshot
User=info
WorkingDirectory=/home/info/fntx-ai-v1/01_backend/alm
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/home/info/fntx-ai-v1/01_backend"
ExecStart=/usr/bin/python3 /home/info/fntx-ai-v1/01_backend/alm/alm_automation.py
StandardOutput=journal
StandardError=journal
Restart=no

[Install]
WantedBy=multi-user.target