# backend/services/fntx-daily-nav.service
[Unit]
Description=FNTX Daily NAV Import Service
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
User=info
Group=info
WorkingDirectory=/home/info/fntx-ai-v1/backend
EnvironmentFile=/home/info/fntx-ai-v1/backend/.env
ExecStart=/home/info/fntx-ai-v1/venv/bin/python /home/info/fntx-ai-v1/backend/services/daily_nav_import.py
StandardOutput=journal
StandardError=journal

# Restart policy
Restart=on-failure
RestartSec=300

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target