# backend/scripts/systemd/fntx-daily-import.timer
[Unit]
Description=Run FNTX Daily Import at 7:00 AM HKT
Requires=fntx-daily-import.service

[Timer]
# Run at 7:00 AM HKT (11:00 PM UTC previous day)
OnCalendar=*-*-* 23:00:00
Persistent=true

[Install]
WantedBy=timers.target