# backend/scripts/systemd/fntx-exercise-detector.timer
[Unit]
Description=Run FNTX Exercise Detection at 7:00 AM HKT
Requires=fntx-exercise-detector.service

[Timer]
# Run at 7:00 AM HKT (11:00 PM UTC previous day)
OnCalendar=*-*-* 23:00:00
Persistent=true

[Install]
WantedBy=timers.target