# backend/services/alm/alm-automation.service
[Unit]
Description=ALM Automation Daily Update
After=network.target postgresql.service

[Service]
Type=oneshot
User=info
WorkingDirectory=/home/info/fntx-ai-v1/backend/alm
Environment="PATH=/home/info/fntx-ai-v1/rl-trading/spy_options/rl_venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/home/info/fntx-ai-v1"
Environment="PYTHONUNBUFFERED=1"

# Load environment variables from .env file
EnvironmentFile=-/home/info/.env

# Run the ALM automation script with virtual environment
ExecStart=/home/info/fntx-ai-v1/rl-trading/spy_options/rl_venv/bin/python3 /home/info/fntx-ai-v1/backend/alm/alm_automation.py

# Restart policy
Restart=on-failure
RestartSec=300

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target